# camera_agent/systemd/sandfest-camera@.service
[Unit]
Description=Texas SandFest camera edge metrics agent (%i)
After=network-online.target time-sync.target
Wants=network-online.target time-sync.target

[Service]
Type=simple
User=sandfest-camera
Group=sandfest-camera
WorkingDirectory=/var/lib/sandfest-camera
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/sandfest
Environment=SANDFEST_CAMERA_ENV=production
EnvironmentFile=/etc/sandfest/camera-%i.env
ExecStartPre=/opt/sandfest/camera_agent/.venv/bin/python -m camera_agent.edge_agent --config /etc/sandfest/camera-agent.json --camera %i --validate-runtime --validate-production
ExecStartPre=/opt/sandfest/camera_agent/.venv/bin/python -m camera_agent.edge_agent --config /etc/sandfest/camera-agent.json --verify-model --model-dir /var/lib/sandfest-camera
ExecStart=/opt/sandfest/camera_agent/.venv/bin/python -m camera_agent.edge_agent --config /etc/sandfest/camera-agent.json --camera %i
Restart=always
RestartSec=5
TimeoutStopSec=20
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/sandfest-camera
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
